# site/abakos-waitlist.service
[Unit]
Description=Abakos waitlist micro-service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
Environment=WAITLIST_PORT=13011
Environment=WAITLIST_FILE=/opt/sites/abakos.ai/data/waitlist.csv
ExecStart=/usr/bin/python3 /opt/sites/abakos.ai/waitlist_service.py
Restart=always
RestartSec=2
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
